# fund.cron
30 6 * * * cd /home/fund-spider && python fund.py