# fund.cron
30 6 * * * python /home/fund-spider/fund.py